# shared/db/backup-db/crontab
# Run every day at 1am
0 1 * * * root /scripts/backup-database.sh
