# shared/db/backup-db/crontab
# Run every day at 1am
0 1 * * * . /etc/environment; /scripts/backup-database.sh >> /var/log/cron/backup.log 2>&1
